# conf/crontab/web-cron
#!/bin/sh
cd $WEBROOT && # run cron cmd
